# voice-assistant.service
[Unit]
Description=Hyprland Voice Assistant (Claude Code)
Documentation=https://github.com/fursman/Assistant
After=pipewire.service
Wants=pipewire.service

[Service]
Type=simple
ExecStart=%h/.local/bin/voice-assistant
Restart=on-failure
RestartSec=5
Environment="PATH=%h/.local/bin:%h/.npm-global/bin:/usr/local/bin:/usr/bin:/bin"
Environment="WAYLAND_DISPLAY=wayland-1"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
# For remote mode, uncomment and set:
# Environment="VOICE_REMOTE=clawbox.local"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voice-assistant

[Install]
WantedBy=default.target
